# ngrokd.service
[Unit]
Description=ngrok server (Go 1.24 version)
After=network.target

[Service]
Type=simple
User=ngrok
Group=ngrok
ExecStart=/opt/ngrok/bin/ngrokd -domain=tunels.tech -httpAddr=:80 -httpsAddr=:443 -tunnelAddr=:4443 -tlsCrt=/opt/ngrok/certs/server.crt -tlsKey=/opt/ngrok/certs/server.key
Restart=always
RestartSec=5
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
